# blue-team/vsftpd-2.3.4-infected/ftpd.service
[Unit]
Description=vsftpd
After=network.target

[Service]
User=root
ExecStart=/usr/local/sbin/vsftpd
Restart=always

[Install]
WantedBy=multi-user.target
